# debian/compatibleone.postrm
#!/bin/sh

case $1 in
	purge)
		update-rc.d compatibleone remove > /dev/null
	;;
esac

#EOF
