# debian/compatibleone.prerm
#/bin/sh

case "$1" in
    remove|deconfigure)
	if [ -x /etc/init.d/compatibleone ]; then
		if [ -x /usr/sbin/invoke-rc.d ]; then
			invoke-rc.d compatibleone stop
		else
			/etc/init.d/compatibleone stop
		fi
	fi
	;;
    upgrade|failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

