# debian/cosacs.postrm
#!/bin/sh

set -e

case $1 in
	purge)
		update-rc.d cosacs remove > /dev/null
	;;
esac

#EOF
